# startup.service
[Unit]
Description=startup

[Service]
ExecStart=/home/pi/startup.sh

[Install]
WantedBy=multi-user.target


